# system/systemd/tiny-dfr.service.d/override.conf
[Unit]
BindsTo=
After=
DefaultDependencies=no
